# Template/etc/inittab
# Start init.d 1st
::sysinit:/etc/init.d/rcS

# Start a couple of shell
tty1::respawn:-/bin/sh
tty2::askfirst:-/bin/sh
tty3::askfirst:-/bin/sh

# Stuff to do when restart init
::restart:/sbin/init

# Stuff to do before reboot
::ctrlaltdel:/bin/reboot
::shutdown:/bin/unmount -a -r

